# systemd/change-wallpaper.service
[Unit]
Description=Change the wallpaper

[Service]
Environment=DISPLAY=:0
ExecStart=/bin/sh %h/.local/scripts/feh_change.sh

